# config/dev/supervisor.conf
[program:code-challenge]
command = /home/vagrant/.local/bin/gunicorn --config /opt/code-challenge/gunicorn_config.py wsgi:app
directory = /opt/code-challenge/repo
user = vagrant

